# crontab.cron
# Look for forecasts every hour.
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin
# Try to get SREF/HREF a few times a day, for our archive (in case the download_and_forecast.rb script borks)
0 6,23 * * * /home/brian/nadocast_operational/get_data/get_sref.rb >> /home/brian/nadocast_operational/get_data/get_sref.log 2>&1
0 6,23 * * * /home/brian/nadocast_operational/get_data/get_href.rb >> /home/brian/nadocast_operational/get_data/get_href.log 2>&1
# Unless specified, the forecast date and time will be inferred from the process start time. Will keep trying to download until the forecast succeeds
# NCEP model status: https://www.nco.ncep.noaa.gov/pmb/nwprod/prodstat/index.html
0  2   * * * source $HOME/.bash_profile >> /home/brian/nadocast_dev/forecaster.log 2>&1 && PUBLISH=true TWEET=true RUN_HOUR=0  HRRR_RAP=false ruby /home/brian/nadocast_dev/lib/download_and_forecast.rb >> /home/brian/nadocast_dev/forecaster.log 2>&1
10 9   * * * source $HOME/.bash_profile >> /home/brian/nadocast_dev/forecaster.log 2>&1 && PUBLISH=true TWEET=true RUN_HOUR=10                ruby /home/brian/nadocast_dev/lib/download_and_forecast.rb >> /home/brian/nadocast_dev/forecaster.log 2>&1
45 13  * * * source $HOME/.bash_profile >> /home/brian/nadocast_dev/forecaster.log 2>&1 && PUBLISH=true TWEET=true RUN_HOUR=14                ruby /home/brian/nadocast_dev/lib/download_and_forecast.rb >> /home/brian/nadocast_dev/forecaster.log 2>&1
# 10 19  * * * source $HOME/.bash_profile >> /home/brian/nadocast_dev/forecaster.log 2>&1 && PUBLISH=true TWEET=true RUN_HOUR=20                ruby /home/brian/nadocast_dev/lib/download_and_forecast.rb >> /home/brian/nadocast_dev/forecaster.log 2>&1
