# crontab.cron
# Look for forecasts every hour.
SHELL=/bin/bash
0 * * * * /home/brian/nadocast/get_data/get_sref.rb >> /home/brian/nadocast/get_data/get_sref.log 2>&1
0 * * * * /home/brian/nadocast/get_data/get_href.rb >> /home/brian/nadocast/get_data/get_href.log 2>&1
15 3 * * * cd /home/brian/nadocast && make forecast_and_publish >> /home/brian/nadocast/forecaster.log 2>&1
30 11 * * * RUN_HOURS=8,9,10 ruby /home/brian/nadocast/get_data/get_rap_current.rb >> /home/brian/nadocast/get_data/get_rap.log 2>&1 && RUN_HOURS=8,9,10 FORECAST_HOURS=1,2,3,4,5,6,7,8,9,10,11,12,13,14,15,16,17,18 ruby /home/brian/nadocast/get_data/get_hrrr.rb >> /home/brian/nadocast/get_data/get_hrrr.log 2>&1 && cd /home/brian/nadocast && make forecast_and_publish >> /home/brian/nadocast/forecaster.log 2>&1
5 16 * * * RUN_HOURS=12,13,14 ruby /home/brian/nadocast/get_data/get_rap_current.rb >> /home/brian/nadocast/get_data/get_rap.log 2>&1 && RUN_HOURS=12,13,14 FORECAST_HOURS=1,2,3,4,5,6,7,8,9,10,11,12,13,14,15,16,17,18 ruby /home/brian/nadocast/get_data/get_hrrr.rb >> /home/brian/nadocast/get_data/get_hrrr.log 2>&1 && cd /home/brian/nadocast && source $HOME/.bash_profile >> /home/brian/nadocast/forecaster.log 2>&1 && make forecast_and_publish >> /home/brian/nadocast/forecaster.log 2>&1
